# deploy/rok1685.service
# systemd unit for running the tracker directly on a VPS (no Docker).
# Install: sudo cp deploy/rok1685.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now rok1685
[Unit]
Description=Rise of Kingdoms 1685 tracker + controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/rok1685
EnvironmentFile=/opt/rok1685/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
